# freeradius/rootfs/etc/s6-overlay/s6-rc.d/freeradius/finish
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Rotina executada quando o FreeRADIUS encerra
# ==============================================================================

bashio::log.warning "Serviço FreeRADIUS encerrado. Código: ${1}"
